# eigennet/files/sbin/link_dbm
#!/bin/sh /etc/rc.common

config_load eigennet

local ip4addr_mesh	; config_get	ip4addr_mesh	network	"ip4addr_mesh"
local iface_mesh=$(ip -4 a s | grep "$ip4addr_mesh" | awk '{print $7}')

[ -e /tmp/link_db ]
{
	rm -rf /tmp/link_db
}

touch /tmp/link_db

iwinfo $iface_mesh a | grep dBm | awk '{ print $1"\t"$2 }' | sed 's/ /\n/g' | sed '/^$/d' > /tmp/link_db

exit 0

